# templates/systemd/plex_repair.service.j2
[Unit]
Description=Automatically repair Plex Media Server DB.
Requires=plex_repair.timer

[Service]
Type=oneshot
ExecStart={{ plex___app._opt_dir ~ "/DBRepair.sh" ~ ' ' ~ plex__repair._options | join(' ') }}
User=root
Group=root
Environment="DBREPAIR_CACHEAGE={{ plex__repair._cache_age }}"
Environment="DBREPAIR_PAGESIZE={{ plex__repair._page_size }}"

[Install]
WantedBy=multi-user.target