# templates/systemd/plex_repair.timer.j2
[Unit]
Description=plex_repair.service timer
Requires=plex_repair.service

[Timer]
OnCalendar={{ plex__repair._schedule }}
Unit=plex_repair.service

[Install]
WantedBy=timers.target