# templates/systemd/plex_repair.service.j2
[Unit]
Description=Automatically repair Plex Media Server DB.
Requires=plex_repair.timer

[Service]
Type=oneshot
ExecStart={{ plex___app._bin_dir ~ '/' ~ plex___repair._package ~ ' ' ~ plex__repair._options | join(' ') }}
WorkingDirectory={{ plex___app._db_dir }}
User=root
Group=root
Environment="DBREPAIR_CACHEAGE={{ plex__repair._cache_age }}"
Environment="DBREPAIR_PAGESIZE={{ plex__repair._page_size }}"

[Install]
WantedBy=multi-user.target