# molecule/default/files/plex_repair.timer
[Unit]
Description=plex_repair.service timer
Requires=plex_repair.service

[Timer]
OnCalendar=weekly
Unit=plex_repair.service

[Install]
WantedBy=timers.target